# cron/nightly_backtest.cron.example
# Example only — no executable logic included
0 23 * * * run nightly backtest and update reports
