# cron/scan_15m.cron.example
# Example only — no executable logic included
*/5 * * * *  scan 15m market-data + signal engine
