# cron/scan_1h.cron.example
# Example only — no executable logic included
*/15 * * * * scan 1H validation context
